# templates/minetest.service.example
[Unit]
Description=MTDESCRIPTION

[Service]
ExecStart=MTEXECUTIONSTRING
Restart=always

[Install]
WantedBy=multi-user.target
